# tools/scripts/reminders/reminderEmails.timer
[Unit]
Description= Run reminderEmails.service every morning

[Timer]
OnCalendar=*-*-* 09:00:00

[Install]
WantedBy=multi-user.target