# tools/scripts/reminders/reminderEmails.service
[Unit]
Description=Reminder emails

[Service]
type=oneshot
ExecStart=/usr/bin/sh -c '~/bin/reminderEmails.sh'